# doc/99-ops/supervisor/hipp.service
# dservice for systemd (CentOS 7.0+)
# by ET-CS (https://github.com/ET-CS)
# "systemctl daemon-reload" reload the config
[Unit]
Description=hipp API
[Service]
Type=forking
ExecStart=supervisord -c /home/supervisor/supervisord.conf
ExecStop=supervisorctl -c /home/supervisor/supervisord.conf  shutdown
ExecReload=supervisorctl -c /home/supervisor/supervisord.conf reload
KillMode=process
Restart=on-failure
RestartSec=42s
[Install]
WantedBy=multi-user.target
